# files/gpslogger.service
[Unit]
After=chrony.service gpsd.service

[Service]
User=pirate
ExecStart=/home/pirate/gpslog.sh

[Install]
WantedBy=default.target
